# config/sddm-theme/etc/systemd/system/sddm-weather.service
[Unit]
Description=Update weather string for SDDM

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-sddm-weather
